# service/seng302-staging-frontend.service
[Unit]
Description=SENG302 Quality Assurance (staging) Frontend service

[Service]
Type=simple
Restart=always
RestartSec=10
User=gitlab-runner
WorkingDirectory=/home/gitlab-runner
ExecStart=/usr/bin/bash staging-frontend.sh

[Install]
WantedBy=multi-user.target
